# backup-files/s3fs.service
[Unit]
After=network.service

[Service]
{TEXT_TO_BE_REPLACED}

[Install]
WantedBy=default.target
